# packaging/systemd/rook-agent.service
[Unit]
Description=RooK servicechannel agent
Documentation=https://github.com/rook
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/rook-agent
ExecStart=/usr/bin/rook-agent service
Restart=on-failure
RestartSec=5s
StateDirectory=rook-agent
RuntimeDirectory=rook-agent
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
